# deploy/vamos-v1-auto-deploy.service
[Unit]
Description=Vamos Subscription Tracker V1 Git auto deploy
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/vamos-subscription-tracker-v1
Environment=APP_DIR=/opt/vamos-subscription-tracker-v1
Environment=BRANCH=develop
Environment=COMPOSE_FILE=docker-compose.v1.yml
Environment=ENV_FILE=.env.v1
ExecStart=/usr/bin/env bash /opt/vamos-subscription-tracker-v1/deploy/git-auto-deploy.sh
